# MESHdesk/files/MESHdesk/hotplug_dir/hotplug.d/iface/70-meshdesk
#!/bin/sh

/etc/MESHdesk/hotplug_lan.lua
