# MESHdesk/files/MESHdesk/hotplug_dir/hotplug.d/iface/80-vpn-routing
#!/bin/sh

/etc/MESHdesk/utils/vpn_routing.lua
